# .config/systemd/user/dotfiles.timer
[Unit]
Description=Run dotfiles.service weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
